# _resources/scripts/usb-monitor.conf
[program:usb-monitor]
directory=/usr/local/bin
environment=
command=python /usr/local/bin/usb-monitor.py
autostart=true
autorestart=true